# mission-control/mission-control-dev.service
[Unit]
Description=Mission Control Next.js app (dev mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=davesalter
WorkingDirectory=/home/davesalter/.openclaw/workspace/mission-control
Environment=NODE_ENV=development
Environment=PORT=3002
ExecStart=/usr/bin/npm run dev -- --hostname 0.0.0.0 --port 3002
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
