# second-brain/second-brain.service
[Unit]
Description=Second Brain Next.js app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=davesalter
WorkingDirectory=/home/davesalter/.openclaw/workspace/second-brain
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start -- --hostname 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
